# .config/systemd/user/taskwarrior_high_council.timer
[Unit]
Description=High Council Taskwarrior Execution

[Timer]
OnCalendar=Mon 10:00
Persistent=false

[Install]
WantedBy=timers.target
