# .config/systemd/user/taskwarrior_virtual_coffee_break.timer
[Unit]
Description=Virtual Coffee Break Taskwarrior Execution

[Timer]
OnCalendar=Fri 09:00
Persistent=false

[Install]
WantedBy=timers.target
