# .config/systemd/user/taskwarrior_weekly.timer
[Unit]
Description=Daily Taskwarrior Execution

[Timer]
OnCalendar=Wed 10:00
Persistent=false

[Install]
WantedBy=timers.target
