# .config/systemd/user/taskwarrior_backend_board_update.timer
[Unit]
Description=Backend Board Update Taskwarrior Execution

[Timer]
OnCalendar=Fri 10:00
Persistent=false

[Install]
WantedBy=timers.target
