# .config/systemd/user/taskwarrior_source_api_chitchat.timer
[Unit]
Description=Source API ChitChat Taskwarrior Execution

[Timer]
OnCalendar=Tue 9:00
Persistent=false

[Install]
WantedBy=timers.target
